# ops/com.madhankumar.jobradar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.madhankumar.jobradar</string>
    <!-- Local LinkedIn-inclusive fetch from the Mac's residential IP (LinkedIn
         blocks datacenter IPs; GitHub Actions runs with --no-linkedin).
         Runs from ~/jobradar-runner, a mirror clone OUTSIDE ~/Documents:
         macOS TCC denies launchd jobs access to Documents/Desktop/Downloads,
         which is why the previous location failed with PermissionError on
         pyvenv.cfg. git pull keeps the mirror current before every run.
         Requires ~/jobradar-runner/pipeline/.env with DATABASE_URL,
         ENCRYPTION_KEY, and ANTHROPIC_API_KEY. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd $HOME/jobradar-runner &amp;&amp; git pull -q origin main; cd pipeline &amp;&amp; .venv/bin/python main.py --no-notion &gt;&gt; /tmp/jobradar.log 2&gt;&amp;1</string>
    </array>
    <!-- Twice daily (local time), interleaved with the Actions runs. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key><string>/tmp/jobradar.out</string>
    <key>StandardErrorPath</key><string>/tmp/jobradar.err</string>
</dict>
</plist>
