# ops/com.madhankumar.jobradar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.madhankumar.jobradar</string>
    <!-- Local LinkedIn-inclusive fetch. Runs the FULL pipeline (ATS + LinkedIn)
         from your Mac's residential IP, since LinkedIn blocks datacenter IPs and
         GitHub Actions skips LinkedIn. Writes to the same Neon database.
         Requires ~/Documents/GitHub/job-radar/pipeline/.env with DATABASE_URL,
         ENCRYPTION_KEY, and ANTHROPIC_API_KEY (see ops/README.md). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>cd $HOME/Documents/GitHub/job-radar/pipeline &amp;&amp; $HOME/Documents/GitHub/job-radar/pipeline/.venv/bin/python main.py --no-notion &gt;&gt; /tmp/jobradar.log 2&gt;&amp;1</string>
    </array>
    <!-- Twice daily (local time), interleaved with the Actions runs. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key><string>/tmp/jobradar.out</string>
    <key>StandardErrorPath</key><string>/tmp/jobradar.err</string>
</dict>
</plist>
