# ops/com.madhankumar.jobradar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.madhankumar.jobradar</string>
    <!-- Local LinkedIn-inclusive fetch. Writes to the same Neon database.
         Requires ~/job-radar/pipeline/.env with DATABASE_URL + ENCRYPTION_KEY. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>cd $HOME/job-radar/pipeline && $HOME/job-radar/pipeline/.venv/bin/python main.py --no-notion >> /tmp/jobradar.log 2>&1</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>12</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key><false/>
    <key>StandardOutPath</key><string>/tmp/jobradar.out</string>
    <key>StandardErrorPath</key><string>/tmp/jobradar.err</string>
</dict>
</plist>
